# package/forked-daapd/files/forked-daapd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2010 OpenWrt.org

BIN=/usr/sbin/forked-daapd
PID=/var/run/forked-daapd.pid
SSD=start-stop-daemon

start() {
	mkdir /tmp/itunes
        $SSD -p $PID -S -x $BIN -- -P $PID
}
	
stop() {
        $SSD -p $PID -K -s SIGKILL
	sleep 1
}
