# package/pot/files/pot.init
#!/bin/sh /etc/rc.common
START=70


pot_start() {
	potd
	potval
}

boot() {
	pot_start
}

pot_stop(){
	killall potval
	killall potd
}

restart(){
	pot_stop
	pot_start
}
