# package/dropbear/files/dropbear.init
#!/bin/sh /etc/rc.common

START=98
STOP=30

DROPBEAR="/usr/sbin/dropbear"
PIDFILE="/var/run/dropbear.pid"
PORT=22

start() {
	# Set TimeZone for console
	tz=$($CONFIG get time_zone)
	echo "$tz" > /tmp/TZ

	# Start SSH server
	$DROPBEAR -p $PORT -a -s -P $PIDFILE
}

stop() {
	# Stop SSH server
	kill `cat $PIDFILE`
}
