# package/irqbalance/files/irqbalance.init
#!/bin/sh /etc/rc.common

START=12
STOP=30

start() {
	/usr/sbin/irqbalance
}

stop() {
	killall -9 irqbalance
}
