# package/igmpproxy/files/igmpproxy.init
#!/bin/sh /etc/rc.common
START=90

boot() {
	/sbin/cmdigmp start
}
