# package/net-scan/files/net-scan.init
#!/bin/sh /etc/rc.common

start()
{
	/usr/sbin/net-scan
}

stop()
{
	killall net-scan
}

boot()
{
	mkdir -p /tmp/netscan
	start
}
