# package/streamboost-misc/files/etc/crontabs/root
# m h  dom mon dow   command
  0 *  *   *   *     logrotate /etc/logrotate.conf -s /var/log/logrotate.status
