# package/trafficmeter/files/traffic_meter.init
#!/bin/sh /etc/rc.common
START=85

start() {
	if [ "$($CONFIG get i_opmode)" = "apmode" -o "$($CONFIG get i_opmode)" = "brmode" ]; then
		echo "traffic_meter do not need start in AP/Bridge mode"
		exit 0
	fi
	/sbin/cmd_traffic_meter start
}

stop() {
	/sbin/cmd_traffic_meter stop
}

boot()
{
	mkdir -p /tmp/traffic_meter
	start
}
