# package/base-files/files/etc/init.d/ddns
#!/bin/sh /etc/rc.common

start() {
	/sbin/cmd_ddns start
}

stop() {
	/sbin/cmd_ddns stop
}

restart() {
	/sbin/cmd_ddns restart
}

