# package/detcable/files/detcable.init
#!/bin/sh /etc/rc.common
START=20

start() {
	# Start detcable
	/usr/bin/detcable 2
}

stop() {
	killall detcable
}

restart() {
	stop
	start
}
