# git_home/wlan-common.git/etc/hotplug.d/wlanlog/00-wlanlog
case "$ACTION" in
    rejected)
            logger "[WLAN access $ACTION: incorrect security] from MAC address $MAC,"
            ;;
esac
