# package/miniupnpd/files/upnp.init
#!/bin/sh /etc/rc.common

start() {
	/sbin/cmdupnp start
}

stop() {
	/sbin/cmdupnp stop
}

restart() {
	/sbin/cmdupnp restart
}

